# higgs/branches/CORBA/code/devices/vaio_tools/Naming_Service.conf
description	"CORBA Naming Service"

start on net-device-up

respawn

exec Naming_Service -ORBEndPoint iiop://higgs2.disam.etsii.upm.es:9876

